# lib/harmony/dist/inputlirc/inputlirc.service
[Unit]
Documentation=man:inputlircd(8)
Description=Zeroconf LIRC daemon using input event devices
After=udev lircd

[Service]
Type=simple
EnvironmentFile=/etc/default/inputlirc
ExecStart=/usr/sbin/inputlircd -f $OPTIONS $EVENTS

[Install]
WantedBy=multi-user.target
